# pisetup/filelist
# List of config files to be modified, appended to or to be replaced by new values
append	bashrc	/home/pi/.bashrc
append	dhcpcd.conf	/etc/dhcpcd.conf
modify	hostapd	/etc/default/hostapd
modify	interfaces	/etc/network/interfaces
modify	sysctl.conf	/etc/sysctl.conf
replace	dnsmasq.conf	/etc/dnsmasq.conf
replace	dnsmasq.service	/etc/systemd/system/dnsmasq.service
replace	forwarder.py	/opt/forwarder/forwarder.py
replace	forwarder.service	/etc/systemd/system/forwarder.service
replace	hostapd.conf	/etc/hostapd/hostapd.conf
replace	statuscli.py	/home/pi/statuscli.py
